# net-p2p/ctcs/files/ctcs.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: ctcs
# REQUIRE: DAEMON
#
# Add the following lines to /etc/rc.conf to run ctcs:
#
# ctcs_enable (bool):	Set it to "YES" to enable ctcs.
#				Default is "NO".
# ctcs_user (user):	Set user to run ctcs.
#				Default is "lastfm".
#

. /etc/rc.subr

name="ctcs"
procname="/usr/bin/perl"
rcvar=${name}_enable
pidfile=/var/run/ctcs.pid

load_rc_config ${name}

: ${ctcs_enable="NO"}
: ${ctcs_user="ctcs"}

command=/usr/bin/su
command_args="-m ${ctcs_user} -c '%%PREFIX%%/bin/ctcs'"

unset ctcs_user

run_rc_command "$1"
